# systemd/daily-saint.service
[Unit]
Description=Daily Orthodox Saint Nostr Post
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/nostr-bot/keys.env
ExecStart=/bin/sh -c '%h/git/personal/daily_saint_bot/.venv/bin/python %h/git/personal/daily_saint_bot/bot.py | %h/git/personal/orthostrbot/.venv/bin/python %h/git/personal/orthostrbot/nostr_post.py'
StandardOutput=journal
StandardError=journal
