# systemd/daily-readings.timer
[Unit]
Description=Daily Scripture Readings Nostr Post Timer

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
