# systemd/daily-saint-random.service
[Unit]
Description=Daily Random Orthodox Saint Nostr Post
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/nostr-bot/keys.env
# Edit the paths below if you cloned the repos somewhere other than ~/git/personal/
ExecStart=/bin/sh -c '%h/git/personal/daily_saint_bot/.venv/bin/python %h/git/personal/daily_saint_bot/bot.py --random | %h/git/personal/orthostrbot/.venv/bin/python %h/git/personal/orthostrbot/nostr_post.py'
StandardOutput=journal
StandardError=journal
